# data/autostart/typio.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Typio
Comment=Start the Typio input method daemon
Exec=typio daemon
TryExec=typio
Icon=typio-keyboard
Terminal=false
StartupNotify=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
Categories=System;Utility;
Keywords=input;ime;keyboard;typio;
